# docker/rootfs/etc/services.d/redis/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

UMASK=${UMASK:-022}
umask ${UMASK}

exec \
    s6-setuidgid root redis-server